# deploy/postern.service
[Unit]
Description=Postern email gateway
Documentation=https://github.com/alexander/postern
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=postern
Group=postern

# Where the binary lives on the host. Adjust if you install elsewhere.
ExecStart=/usr/local/bin/postern

# Persistent state. systemd creates /var/lib/postern owned by the service user.
StateDirectory=postern
WorkingDirectory=/var/lib/postern

# Required: 32-byte master key (hex or base64). KEEP THIS BACKED UP — losing
# it means any encrypted SMTP password in the DB becomes unrecoverable.
EnvironmentFile=/etc/postern/postern.env

Environment=POSTERN_DB_PATH=/var/lib/postern/postern.db
Environment=POSTERN_LISTEN_ADDR=127.0.0.1:8080

Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
ReadWritePaths=/var/lib/postern
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
